# deploy/cron/lansia-papua-cleanup-exports.cron.example
# Laravel scheduler already runs exports:prune daily. Use this only when not using schedule:run.
0 3 * * * cd /var/www/lansia-papua && php artisan exports:prune --delete >> storage/logs/export-cleanup.log 2>&1
